# macos/Users/sametaor/Library/LaunchAgents/com.github.mpd-now-playable.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.github.mpd-now-playable</string>
	<key>Program</key>
	<string>/Users/sametaor/.local/bin/mpd-now-playable</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/mpd-now-playable.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/mpd-now-playable.log</string>
</dict>
</plist>
